# DEBIAN/prerm
#!/bin/bash
systemctl stop leek-duck
systemctl disable leek-duck
exit 0